# .config/systemd/user/work-vpn.service
[Unit]
Description=Work VPN

[Service]
Type=idle
ExecStart=/home/rekyuu/.local/bin/work-vpn
KillSignal=SIGTERM

[Install]
WantedBy=default.target